# docker/crontabs/build
#!/bin/bash
docker build -t zt1115798334/crontabs .